# systemd/audioServer.service
[Unit]
Description=Audio Server Service
After=mpv.service
Requires=mpv.service

[Install]
WantedBy=hostapd.service
# WantedBy=multi-user.target

[Service]
ExecStart=/usr/bin/audioServer 0.0.0.0 443 /var/audioserver

Restart=always
RestartSec=2

